# deploy/roles/supervisor/files/supervisord.conf
description     "supervisord"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/supervisor/venv/bin/supervisord --nodaemon --configuration /opt/supervisor/supervisord.conf
